# monitor/monitor.service
[Unit]
Description=My watchlog service
[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/monitor
ExecStart=/opt/monitor.sh $WORD $LOG
